# data/supervisor/aitodo.conf
[program:aitodo-api]
command=/opt/conda/envs/aitodo/bin/uvicorn api:app --host 0.0.0.0 --port 6100 --workers 8 --timeout-keep-alive 300 --log-level info
directory=/aitodo/api
autorestart=true
autostart=true
stopasgroup=true
killasgroup=true
environment=ES_CONNECT_SETTING_HOST="http://localhost:9200",EMAIL_SENDER="",EMAIL_AUTH="",EMAIL_SMTP="",EMAIL_PORT="",EMAIL_LABEL=""

[program:aitodo-consumer]
command=/opt/conda/envs/aitodo/bin/python consumer.py -n 4
directory=/aitodo/api
autorestart=true
autostart=true
stopasgroup=true
killasgroup=true
environment=ES_CONNECT_SETTING_HOST="http://localhost:9200",EMAIL_SENDER="",EMAIL_AUTH="",EMAIL_SMTP="",EMAIL_PORT="",EMAIL_LABEL=""


[program:aitodo-client]
command=npm run dev
directory=/aitodo/client
autorestart=true
autostart=true
stopasgroup=true
killasgroup=true

