# data/supervisor/service.conf
[supervisord]
nodaemon=true

[inet_http_server]         
port=0.0.0.0:9001        
username=admin
password=111111

[program:redis]
command=/usr/bin/redis-server --requirepass "111111"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=ES_JAVA_OPTS="-Xms512m -Xmx512m"

[program:coder]
command=bash /bin/code-server
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log


